# templates/launchagent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  dev-rules cross-machine sync agent.

  Why this exists:
    The single source of truth for rules lives at ~/Codes/dev-rules/. When
    someone (or another machine) pushes to its configured origin, this machine's
    local mirror — and therefore the symlinks under ~/.cursor/rules and ~/.claude/
    — silently fall behind unless something pulls. That "something" is this
    LaunchAgent.

  What it does:
    Every 30 minutes, runs `dev-rules/sync.sh --pull`, which:
      1. fetches origin and fast-forwards ~/Codes/dev-rules
      2. refreshes home symlinks (idempotent)
      3. fans out real-copy rules to every registered consumer project

  Pair with `sync.sh --push` (the proactive path, run by the editor): together
  they cover both the same-machine and cross-machine update modes.

  Installation: do NOT edit this file by hand. Run install-launchagent.sh which
  templates {{HOME_CANONICAL}} and registers the agent with launchctl.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>local.dev-rules.sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>{{HOME_CANONICAL}}/sync.sh --pull</string>
    </array>

    <key>StartInterval</key>
    <integer>1800</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>{{HOME}}/Library/Logs/dev-rules-sync.log</string>

    <key>StandardErrorPath</key>
    <string>{{HOME}}/Library/Logs/dev-rules-sync.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>{{HOME}}</string>
    </dict>
</dict>
</plist>
